# systemd/nmct-box-flask.service
[Unit]
Description=uWSGI instance to serve student uploads
After=network.target

[Service]
User=nmct
Group=www-data
WorkingDirectory=${NMCT_HOME}/src
Environment="PATH=${NMCT_HOME}/env/bin"
Environment="NMCT_HOME=${NMCT_HOME}"
ExecStart="${NMCT_HOME}/env/bin/uwsgi" --ini "${NMCT_HOME}/src/nmct/web/wsgi-student.ini"

[Install]
WantedBy=multi-user.target
